# playbooks/templates/bdd-cron.timer.j2
[Unit]
Description=bdd check verify: {{ verify_data }} 

[Timer]
OnCalendar={{ timer_pattern }}

[Install]
WantedBy=timers.target